# pkgbuilds/rpi4-network/rpi4-network.install
post_install() {
        systemctl enable systemd-networkd.service
        systemctl enable systemd-resolved.service
        ln -sf /run/systemd/resolve/stub-resolv.conf /etc/resolv.conf
}
