# pkgbuilds/rpi4-bluetooth/rpi4-bluetooth.install
post_install() {
  echo 'blacklist btsdio' >> /etc/modprobe.d/raspi-blacklist.conf
  systemctl enable bluetooth.service
  systemctl enable brcm43438.service
}
