# pkgbuilds/rpi4-firstboot/10-systemd-firstboot.service
[Service]
ExecStart=
ExecStart=/usr/bin/systemd-firstboot --prompt
ExecStartPost=/usr/bin/systemd-firstboot-rpi4

[Install]
WantedBy=
WantedBy=sysinit.target
